# etc/scripts/supervisor/publisher-worker.conf
[program:messenger-consume]
directory=/var/www/publisher/
command=php bin/console messenger:consume async_content_push async_analytics_event async_webhooks async_analytics_export async_image_conversion --time-limit=3600
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/publisher/var/log/stdout.log
stderr_logfile=/var/www/publisher/var/log/stderr.log
