# etc/scripts/supervisor/supervisor_analytics_export.conf
[program:publisher_analytics_export]
directory=/var/www/publisher/
command=php bin/console messenger:consume async --env=prod
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/www/publisher/var/log/stdout.log
stderr_logfile=/var/www/publisher/var/log/stderr.log
