# etc/scripts/supervisor/supervisor-websocket.conf
[program:publisher-websocket]
directory=/var/www/publisher/
command=php bin/console gos:websocket:server --env=prod --time-limit=3600
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/publisher/var/log/stdout.log
stderr_logfile=/var/www/publisher/var/log/stderr.log
