# scripts/supervisor/supervisor_content_push.conf
[program:publisher_analytics]
directory=/var/www/publisher/
command=php app/console rabbitmq:consumer content_push --env=prod
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/www/publisher/app/logs/supervisor_stdout.log
stderr_logfile=/var/www/publisher/app/logs/supervisor_stderr.log
