# scripts/supervisor/supervisor_analytics.conf
[program:publisher_analytics]
directory=/var/www/publisher/
command=php app/console rabbitmq:consumer analytics_event --messages=250 --env=prod
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/www/publisher/var/logs/supervisor_stdout.log
stderr_logfile=/var/www/publisher/var/logs/supervisor_stderr.log
