# scripts/supervisor/supervisor_webhooks.conf
[program:publisher_webhooks]
directory=/var/www/publisher/
command=php app/console rabbitmq:consumer send_webhook --messages=250 --env=prod
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/www/publisher/var/log/stdout.log
stderr_logfile=/var/www/publisher/var/log/stderr.log
