# scripts/supervisor/supervisor_commands.conf
[program:run_command]
directory=/var/www/publisher/
command=php app/console rabbitmq:consumer run_command --env=prod
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/www/publisher/app/logs/stdout.log
stderr_logfile=/var/www/publisher/app/logs/stderr.log
