# supervisor.ini
[supervisord]
nodaemon=true

[program:websocket]
command = php /var/www/publisher/bin/console gos:websocket:server
user = root
autostart = true
redirect_stderr=true

[program:content_push]	
command = php /var/www/publisher/bin/console messenger:consume async_content_push
user = root
user = root	
autostart = true	
autorestart=true
stopsignal=INT	
stopwaitsecs=60	

[program:webhooks]
command = php /var/www/publisher/bin/console messenger:consume async_content_push
user = root
autostart = true
autorestart=true
