# supervisor.ini
[supervisord]
nodaemon=true

[program:websocket]
command = php /var/www/publisher/bin/console gos:websocket:server
user = root
autostart = true
autorestart=true

[program:analytics]
command = php /var/www/publisher/bin/console messenger:consume async_analytics_event
user = root
autostart = true
autorestart=true

[program:content_push]
command = php /var/www/publisher/bin/console messenger:consume async_content_push
user = root
autostart = true
autorestart=true

[program:webhooks]
command = php /var/www/publisher/bin/console messenger:consume async_webhooks
user = root
autostart = true
autorestart=true
