# supervisor.ini
[supervisord]
nodaemon=true

[program:websocket]
command = php /var/www/publisher/bin/console gos:websocket:server
user = root
autostart = true
redirect_stderr=true


[program:content_push]
command = php /var/www/publisher/bin/console rabbitmq:consumer content_push
user = root
autostart = true
redirect_stderr=true
stdout_logfile=/var/www/publisher/app/logs/supervisor_stdout.log
stderr_logfile=/var/www/publisher/app/logs/supervisor_stderr.log
