# etc/docker/supervisor/supervisor.ini
[supervisord]
nodaemon=true

[program:websocket]
command = php /var/www/publisher/bin/console gos:websocket:server
user = root
autostart = true

[program:analytics]
command = php /var/www/publisher/bin/console rabbitmq:consumer analytics_event --messages=250
user = root
autostart = true

[program:content_push]
command = php /var/www/publisher/bin/console rabbitmq:consumer content_push
user = root
autostart = true

[program:webhooks]
command = php /var/www/publisher/bin/console rabbitmq:consumer send_webhook --messages=250
user = root
autostart = true
