# etc/docker/supervisor/websocket.conf
[supervisord]
nodaemon=true

[program:websocket]
command = php /var/www/publisher/bin/console gos:websocket:server
user = root
autostart = true
redirect_stderr=true